# routers/alcor/rc.local
# Put your custom commands here that should be executed once
# the system init finished. By default this file does nothing.

# Reload contract module to prevent it get loaded before network init
rmmod nf_conntrack_netlink
modprobe nf_conntrack_netlink

exit 0
