# routers/qiaoying/uci-defaults/90-dhcp-vendorid
#!/bin/sh

sed -i 's/\${vendorid:+-V "\$vendorid"}/\${vendorid:+-V "" "-x 0x3c:\$vendorid"}/' /lib/netifd/proto/dhcp.sh

exit 0